# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:scratch]
command = /home/docker/code/ionic-theme
stdout_syslog=true
stderr_syslog=true
